# systemd/subgen.service
[Unit]
Description=subgen mihomo subscription server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The service runs from the git checkout in the server user's home. deploy.sh
# builds ./subgen here and restarts the unit. Adjust the paths below if your
# checkout lives elsewhere than /home/server/subgen.
WorkingDirectory=/home/server/subgen
EnvironmentFile=/home/server/subgen/.env
Environment=GOMEMLIMIT=80MiB GOMAXPROCS=1
ExecStart=/home/server/subgen/subgen
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3

# Hardening. The store is written under the checkout's db/ dir; the TLS cert is
# read-only. Everything else is read-only.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/root/cert
ReadWritePaths=/home/server/subgen/db
PrivateTmp=true
MemoryMax=100M
CPUQuota=40%

[Install]
WantedBy=multi-user.target
